# platform/web-git/docker/s6/gogs/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=git
exec gosu $USER /app/boxlinker-git/boxlinker-git web
