# platform/other/livenowhy/Dockerfile/Do-old/python/2.7.10/flask-ubuntu/app.conf
[program:app]
command=/usr/bin/gunicorn -w 1 wsgiapp:application
directory=/srv/www
user=www-data